# sys-container/dockerfiles/jenkins-syscont/supervisord.conf
[supervisord]
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

[program:dockerd]
command=/usr/bin/dockerd
priority=1
autostart=true
autorestart=true
startsecs=0

[program:jenkins]
command=/usr/local/bin/jenkins.sh
priority=2
autostart=true
autorestart=true
startsecs=0
username=jenkins